# code/backend/deploy/hama-fastapi.service.example
[Unit]
Description=Hama FastAPI backend
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/main/code/backend/src/main/python
EnvironmentFile=/home/ubuntu/main/code/backend/src/main/python/.env
Environment=HAMA_OPENSEARCH_ENABLED=true
ExecStart=/home/ubuntu/main/code/backend/src/main/python/.venv/bin/uvicorn api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
